# etc/systemd/system/swayosd-libinput-backend.service.d/override.conf
[Unit]
PartOf=
PartOf=default.target
After=
After=default.target

[Install]
WantedBy=
WantedBy=default.target

