# etc/systemd/system/NetworkManager-wait-online.service.d/override.conf
[Service]
; ExecStart=
; ExecStart=/usr/bin/nm-online -q
Environment=NM_ONLINE_TIMEOUT=30

