# config/systemd/user/btrfs-space-monitor.timer
[Unit]
Description=Run Btrfs space monitor every 30 minutes
Requires=btrfs-space-monitor.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target