# config/systemd/user/screenshot-cleanup.timer
[Unit]
Description=Weekly screenshot cleanup

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
