# etc/systemd/system/btrbk-usb.timer
[Unit]
Description=btrbk USB backup (every 4 hours)
Documentation=man:btrbk(1)

[Timer]
OnBootSec=15min
OnUnitActiveSec=4h
Persistent=false

[Install]
WantedBy=timers.target
