# etc/systemd/system/containerd-cleanup.timer
[Unit]
Description=Weekly Containerd Storage Cleanup
Documentation=man:containerd(8)

[Timer]
# Run every Sunday at 3 AM
OnCalendar=Sun *-*-* 03:00:00
# Run immediately if the system was off at the scheduled time
Persistent=true
# Randomize start time by up to 30 minutes to avoid load spikes
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target