# etc/systemd/system/mnt-usb_backups.automount
[Unit]
Description=USB encrypted btrfs backup automount
After=dev-mapper-usb_backups.device
Requires=dev-mapper-usb_backups.device

[Automount]
Where=/mnt/usb_backups
TimeoutIdleSec=14400

# [Install] removed — drive is offline; units can still be started manually
# WantedBy=multi-user.target