# etc/systemd/system/notify-success@.service
[Unit]
Description=Success notification for %i

[Service]
Type=oneshot
ExecStart=/usr/bin/gotify push -t 'Unit Completed' '%i'
